# backup-system/systemd/win11-trading-backup.service
[Unit]
Description=Win11-Trading powered-off VM backup
Documentation=file://__BACKUP_ROOT__/README.md
After=virtqemud.service libvirtd.service

[Service]
Type=oneshot
ExecStart=__BACKUP_ROOT__/bin/backup-vm.sh
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
