# backup-system/systemd/win11-trading-backup.timer
[Unit]
Description=Run Win11-Trading VM backup weekly when powered off

[Timer]
OnCalendar=Sat *-*-* 14:00:00
Persistent=true
RandomizedDelaySec=5m

[Install]
WantedBy=timers.target
